# modules/autostart/launchd/com.eva.deal-scout.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.eva.deal-scout</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/vineetkumar/Eva/modules/autostart/run-service.sh</string>
        <string>deal-scout</string>
        <string>main.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/vineetkumar/Eva/modules/deal-scout</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>ProcessType</key>
    <string>Background</string>
    <key>StandardOutPath</key>
    <string>/Users/vineetkumar/Eva/logs/eva-deal-scout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/vineetkumar/Eva/logs/eva-deal-scout.error.log</string>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/vineetkumar</string>
    </dict>
</dict>
</plist>
